# primary_io_manager/primary_io_manager.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Primary Io Management daemon"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
expect fork
# primary_io_manager maintains some state but can otherwise recover.
oom score -100
respawn
# if primary_io_manager has a non-error exit, don't attempt to respawn.
normal exit 0

exec minijail0 \
    --config /usr/share/minijail/primary_io_manager.conf -- \
    /usr/bin/primary_io_manager
