# login_manager/init/upstart/arc-manager.conf
# Copyright 2025 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS user interface"
author        "chromium-os-dev@chromium.org"

start on starting ui
stop on stopping ui

# OK to kill if absolutely needed.
oom score -100

# Android guest os sets rtprio via setrlimit. Unless we set this to arc-manager
# the init process would require CAP_SYS_RESOURCE in the container, which is
# prohibited.
limit rtprio 10 10

exec arc_manager
