# login_manager/init/upstart/devicesettings-migration.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Point /var/lib/devicesettings to /var/lib/whitelist"
author          "chromium-os-dev@chromium.org"

start on starting boot-services

# Creates a bind mount of /var/lib/devicesettings -> /var/lib/whitelist to
# allow code to migrate from /var/lib/whitelist to /var/lib/devicesettings.

task
exec mount --bind /var/lib/whitelist /var/lib/devicesettings