# arc/container/file-syncer/init/arc-file-syncer.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run file syncer for ARC++"
author        "chromium-os-dev@chromium.org"

start on start-arc-instance
stop on stop-arc-instance or stopped ui

# 'never' since terminating this job after it starts can cause ARC boot failure.
oom score never

script
  {
    echo "Start arc-file-syncer"
    set -x

    # Start constructing minijail0 args...
    args="minijail0"

    args="$args -- /usr/bin/arc-file-syncer"

    logger -t "${UPSTART_JOB}" "Executing: $args"

    exec $args
  } 2>&1 | logger -t "${UPSTART_JOB}"
end script
