# ml_core/cacher/init/ml-core-cacher.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Compilation cacher for ML models used by ml-core"
author "chromium-os-dev@chromium.org"

start on start-user-session and stopped imageloader-init
stop on stopping ui

oom score -100

expect fork

# This service will load the ML models and processing shaders, and run a
# tuning process to compile and optimize them for that specific device
# GPU/NPU. This happens on the GPU/NPU itself, and so the sandbox needs
# access to those devices.

exec minijail0 --config /usr/share/minijail/ml-core-cacher.conf \
    -- /usr/bin/ml_core_cacher
