# ml_core/cacher/init/ml-core-cacher.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Compilation cacher for ML models used by ml-core"
author "chromium-os-dev@chromium.org"

start on start-user-session and started system-services and stopped imageloader-init
stop on stopping ui

oom score -100

expect fork

# This service will load the ML models and processing shaders, and run a
# tuning process to compile and optimize them for that specific device
# GPU/NPU. This happens on the GPU/NPU itself, and so the sandbox needs
# access to those devices.

script
  # Addiotinal conditional args here, because minijail-config-file doesn't
  # support conditional args.
  EXTRA_ARGS=""

  # Need access to /etc/openvino to read OpenVINO config file for Intel NPU.
  # The path only exists on NPU-enabled devices, so we mount it conditionally.
  if [ -d /etc/openvino ]; then
    EXTRA_ARGS="$EXTRA_ARGS -b /etc/openvino"
  fi

  # Need access to /etc/ml_core to read tflite stable delegate settings file.
  if [ -d /etc/ml_core ]; then
    EXTRA_ARGS="$EXTRA_ARGS -b /etc/ml_core"
  fi

  # Need access to /var/lib/ml_core/openvino_cache to write cache back to disk.
  if [ -d /var/lib/ml_core/openvino_cache ]; then
    EXTRA_ARGS="$EXTRA_ARGS -b /var/lib/ml_core/openvino_cache,,1"
  fi

  exec minijail0 --config /usr/share/minijail/ml-core-cacher.conf \
      ${EXTRA_ARGS} \
      -- /usr/bin/ml_core_cacher
end script
