# featured/share/featured.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS feature daemon"
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn

# This service enables kernel features. It can be killed,
# at the cost of some features not being enabled.
oom score -100

exec /sbin/featured

