# featured/share/featured.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS feature daemon"
author          "chromium-os-dev@chromium.org"

# Chrome does not (yet, as of 2021-12) make finch state available until the
# first restart of the browser -- that is, when a user logs in for the first
# time.
start on start-user-session
stop on shutdown
stop on reboot
respawn

# This service enables kernel features. It can be killed,
# at the cost of some features not being enabled.
oom score -100

exec /sbin/featured
