# imageloader/init/widevine-cdm.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Mount component updated Widevine CDM if it exists"
author "jrummell@chromium.org"

# Note that "starting ui" blocks UI startup. This is intentional here as the
# Widevine CDM component must be available when the browser starts since the
# sandbox/zygote can currently only access binaries if they're available at
# startup.
start on starting ui or started imageloader-init
task

oom score never

# Try to mount Widevine CDM if it exists.
# If it does not exist, that is OK: just let imageloader fail.
# The image is mounted in a pre-determined location that Chrome
# (chrome/browser/component_updater/widevine_cdm_component_installer.cc)
# knows about.
exec /usr/sbin/imageloader --mount --mount_component=WidevineCdm
