# hps/daemon/init/hpsd.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start Chrome OS HPS service"
author        "chromium-os-dev@chromium.org"

# Environment variables for daemon configuration.
# They are empty here and should be overridden in hpsd.override
env HPS_HW_FLAGS=
env FW_APPL_PATH=
env FW_SPI_PATH=
env FW_VERSION=0
#
# Development options for prototyping
env UART_DEV=

# Start when mems_setup has set proper group and ownership.
start on started boot-services
stop on stopping boot-services
expect fork
respawn
respawn limit 5 30

# Make hpsd killable, because if it has a leak it's better to
# restart it than to OOM-panic.
oom score -100

# Minijail settings:
# -i exit immediately after fork
# -u run as user hpsd
# -g run as group hpsd
# -n set no_new_privs
# -R RLIMIT_NICE,40,40 to allow niceness of -20
#
# Namespaces:
# -N enter new cgroup namespace
# --uts enter new UTS/hostname namespace
# -e enter new network namespace
# -p enter new pid namespace
# -P pivot root to /mnt/empty
#
# Mounts:
# -k to mount tmpfs at /run and /var (writable)
# -b /usr/sbin/hpsd to make hpsd accessible
# -b /usr/lib/firmware/hps to make firmware images readable
# -b /var/lib/metrics (writable) to enable UMA
# -r to remount /proc readonly
#
# For I2C:
# -b /dev (writable)
# -b /sys
# -b /sys/bus
# -b /sys/class
# -b /sys/devices (writable)
#
# For DBUS:
# -b /run/dbus
exec minijail0 -i -u hpsd -g hpsd \
  -n -R RLIMIT_NICE,40,40 \
  -N --uts -e -p -P /mnt/empty \
  -k 'tmpfs,/run,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC' \
  -k 'tmpfs,/var,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC' \
  -b /usr/sbin/hpsd -b /usr/lib/firmware/hps \
  -b /var/lib/metrics,,1 -r \
  -b /dev,,1 -b /sys -b /sys/bus -b /sys/class -b /sys/devices,,1 \
  -b /run/dbus \
  -- /usr/sbin/hpsd --skipboot --test --version="${FW_VERSION}" \
     --mcu_path="${FW_APPL_PATH}" --spi_path="${FW_SPI_PATH}" ${HPS_HW_FLAGS}