# hps/daemon/init/hpsd_roflash.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Flash ro of hps before hpsd starts"
author        "chromium-os-dev@chromium.org"

# hpsd will still start if this fails
# we could add this to hpsd.conf if needed:
# start on stopped hpsd_roflash RESULT=ok
start on starting hpsd
task

# This must not be killed as there is a small chance
# it could leave the HPS device bricked.
oom score never
# Typical vsize for hps-factory is 2756 KiB, allow some headroom
limit as 20000000 unlimited

env GPIO_CHIP=0
env GPIO_LINE=327

# Initialise the HPS before starting hpsd:
# 1. Power-cycle the HPS peripheral to ensure it is in a good state.
# 2. If HPS is in the STM bootloader (I2C address 0x51) because its flash is
#    blank, write stage0. If HPS is already in stage0 (I2C address 0x30)
#    this step does this nothing.
# 3. Power-cycle the HPS peripheral again. This is only necessary if we just
#    wrote a new stage0, but it's simpler to unconditionally do it again.
script
  gpioset "${GPIO_CHIP}" "${GPIO_LINE}=0" \
  && sleep 0.1 \
  && gpioset "${GPIO_CHIP}" "${GPIO_LINE}=1" \
  && sleep 1 \
  && syslog-cat --identifier="${UPSTART_JOB}" -- \
    hps-factory --dev=/dev/i2c-hps-controller \
    write-stage0 --existing-stage0-ok \
  && gpioset "${GPIO_CHIP}" "${GPIO_LINE}=0" \
  && sleep 0.1 \
  && gpioset "${GPIO_CHIP}" "${GPIO_LINE}=1" \
  && sleep 1
end script
