# flex_hwis/init/upstart/flex_boot_log.conf
# Copyright 2025 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Copy the crdyboot log from a UEFI variable"
author          "chromium-os-dev@chromium.org"

start on started system-services

oom score -100
task

env efi_dir="/sys/firmware/efi"
env var_name="CrdybootLog"
env var_guid="2a6f93c9-29ea-46bf-b618-271b63baacf3"

pre-start script
  # If `efi_dir` does not exist, then stop the job since the system was
  # not booted via UEFI.
  if [ ! -e "${efi_dir}" ]; then
    stop
    exit 0
  fi
end script

# Copy the crdyboot log [1] from a UEFI variable to /var/log.
#
# The first four bytes are skipped because efivarfs uses them to encode
# variable attributes [2].
#
# [1]: https://source.chromium.org/chromiumos/chromiumos/codesearch/+/HEAD:src/platform/crdyboot/docs/logging.md
# [2]: https://docs.kernel.org/filesystems/efivarfs.html
exec dd bs=1 skip=4 \
  if="${efi_dir}/efivars/${var_name}-${var_guid}" \
  of=/var/log/crdyboot.log
