# flex_hwis/init/upstart/flex_hwis.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run hwis_tool on ChromeOS Flex devices"
author          "chromium-os-dev@chromium.org"

start on start-user-session
oom score -100
tmpfiles /usr/lib/tmpfiles.d/flex_hwis.conf

task

exec minijail0 --config /usr/share/minijail/flex_hwis.conf -- \
  /usr/bin/flex_hwis_tool
