# flex_hwis/init/upstart/flex_hardware_cache.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Cache hardware data on ChromeOS Flex"
author          "chromium-os-dev@chromium.org"

start on started cros_healthd
oom score -100
tmpfiles /usr/lib/tmpfiles.d/flex_hardware_cache.conf

task

exec minijail0 --config /usr/share/minijail/flex_hardware_cache.conf -- \
    /usr/bin/flex_hwis_tool --cache
