# flex_hwis/init/upstart/flex_disk_metrics.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run flex_disk_metrics on ChromeOS Flex devices"
author          "chromium-os-dev@chromium.org"

start on started system-services

oom score -100
task

exec minijail0 --config /usr/share/minijail/flex_disk_metrics.conf -- \
  /usr/bin/flex_disk_metrics
