# farfetchd/init/farfetchd.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Farfetchd daemon for ChromeOS"
author          "chromium-os-dev@chromium.org"

# Start once d-bus is ready.
start on started system-services
stop on stopping system-services

# If the job respawns 10 times in 10 seconds, stop trying.
respawn
respawn limit 10 10

# Service can recover on restart
oom score -100

exec /usr/sbin/farfetchd
