# odml/init/odmld.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS device On-Device ML service."
author          "chromium-os-dev@chromium.org"

start on started system-services and started dlcservice and stopped imageloader-init
stop on stopping system-services
respawn

# Fix glibc malloc's trim_threshold tunable to its default, disabling automatic
# adjustment of this tunable. This helps greatly with unusual allocation
# patterns seen in mediapipe.
env GLIBC_TUNABLES=glibc.malloc.trim_threshold=131072  # 128*1024 bytes

tmpfiles /usr/lib/tmpfiles.d/odml.conf

oom score -100

# The on device model service with APU backend will open a lot of files so the
# process needs a higher limit.
limit nofile 16384 16384

expect fork

pre-start script
  # TODO(b/350812378): Use a proper feature name.
  if ! [ "$(feature_check --feature_level)" -ge "1" ]; then
    logger -t "${UPSTART_JOB}" "Not available on this device."
    stop
    exit 0
  fi

  # TODO(b/365876446): Remove after the auto probe finished.
  modprobe mtk_apusys || true
end script

script
  set --
  # Handle optional bind mounts.

  # ARM GPU Devices.
  if [ -c "/dev/mali0" ]; then
    set -- "$@" -b /dev/mali0
  fi

  # MTK APU devices.
  if [ -c "/dev/apusys" ]; then
    set -- "$@" -b /dev/apusys
  fi

  if [ -d "/dev/dma_heap" ]; then
    set -- "$@" -b /dev/dma_heap
  fi

  exec minijail0 --config /usr/share/minijail/odmld.conf \
  "$@"                                                   \
  -- /usr/bin/odmld
end script
