# gsclog/init/gsclog.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "gsc log concatenator"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

normal exit 0
respawn

oom score -100

exec periodic_scheduler --period=300 --timeout=30 --task_name=gsclog \
  --start_immediately -- \
  /usr/bin/gsclog
