# discod/init/discod.conf
# Copyright 2022 The ChromiumOS Authors.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the disc control daemon"
author          "chromium-os-dev@chromium.org"

oom score -100

start on started boot-services
stop on stopping boot-services

respawn
expect fork

exec discod
