# vm_tools/reference_vm/data/usr/local/lib/systemd/system/install-refvm.service
[Unit]
Description=Reference VM installer

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/install-refvm
StandardOutput=journal+console
StandardError=journal+console
