# minios/ramfs/etc/init/minios.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Startup for MiniOS script"
author          "chromeos-core-services@google.com"

start on started dbus and started frecon and started shill and started udev
expect fork
respawn

oom score never

pre-start script
  backlight_device=
  brightness=
  for backlight_device in /sys/class/backlight/*; do

    # In order to set a backlight to maximum brightness, we do the following:
    # - Turn the backlight device on (using bl_power).
    # - Read max_brightness and pipe that into the desired brightness.
    # - Verify that the backlight's actual brightness matches our expectation.
    #
    # The hardware is particularly volatile during startup due to modeset and
    # hw initialization. We want to give ourselves the best chance of success,
    # so we'll try up to 5 times before continuing with recovery.
    for _ in $(seq 1 5); do
      echo 0 > "${backlight_device}/bl_power" # 0 = FB_BLANK_UNBLANK (ie: on)
      brightness=$(cat "${backlight_device}/max_brightness")
      echo "${brightness}" > "${backlight_device}/brightness"
      if [ "$(cat "${backlight_device}/actual_brightness")" = \
           "${brightness}" ]; then
        break
      fi
    done
  done
end script

exec /usr/bin/minios
