# arc/setup/init/arcpp-post-login-services.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Job to start ARC++ services after user login"
author        "chromium-os-dev@chromium.org"

start on continue-arc-boot
stop on stop-arc-instance or stopping ui

# The following environment variables are passed from session_manager
# and are imported from the event that starts the job.
import CHROMEOS_USER
import CONTAINER_PID

export CHROMEOS_USER
export CONTAINER_PID
