# vm_tools/init/vm_concierge_nodefer.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Check an config value and start vm_concierge if it's disabled"
author         "chromium-os-dev@chromium.org"

# This job just provides a trigger event for vm_concierge to start at
# "start-user-session" for a finch experiment. vm_concierge.conf waits for
# "started vm_concierge_nodefer" event or "started-user-sesison" event by the
# "start on" stanza. This job starts successfully if the DeferConciergeStartup
# config is disabled (or failed to check the config value as the default
# behavior).
# TODO(b/353431869): remove this conf when the finch experiment is finished

start on start-user-session
stop on started vm_concierge

# This job itself does nothing.
oom score 1000

pre-start script
  # Check if the config to defer concierge startup is enabled.
  if /bin/sh /usr/share/cros/init/check_defer_concierge_config.sh ; then
    # The config is enabled. Stop this job to not start concierge.
    stop
    exit 0
  fi
end script
