# init/upstart/test-init/factory-iptables.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS factory post hook after patchpanel"
author        "chromium-os-dev@chromium.org"

# For factory software to expose Goofy rpc and overlord lan discovering, we add
# addtional rules in iptables.

oom score never

# After CL:5311860, `start patchpanel` will restore iptables to a state
# regardless of the current contents of iptables. Thus, we configure the
# iptables with `start on started patchpanel`.
start on started patchpanel
task

script
  . /usr/share/cros/factory_utils.sh

  if ! is_factory_mode; then
    exit 0
  fi

  /usr/local/factory/init/startup iptables
end script
