# init/upstart/test-init/report-initstate-to-ec.conf
# Copyright 2025 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report OS init state to EC"
author        "chromium-os-dev@chromium.org"

oom score -100

# These events are specifically chosen to debug unreachable lab devices.
# These events are not overlapping and hence should be safe to report from
# an upstart script.
start on starting boot-complete or starting network-services or shill-connected

task

script
  # This script is installed only in test images and hence will run only in
  # developer and lab machines.
  logger -t "${UPSTART_JOB}" "${UPSTART_EVENTS} ${JOB}"
  ectool console_print "${UPSTART_EVENTS} ${JOB}"
  if [ ${UPSTART_EVENTS} = "shill-connected" ]; then
    sleep 30
    ip_address=$(ifconfig eth0 | grep "inet " | tr -s ' ' | cut -d ' ' -f3)
    logger -t "${UPSTART_JOB}" "DUT IP Address: ${ip_address}"
    ectool console_print "${ip_address}"
  fi
end script
