# libhwsec-foundation/init/hwsec-shutdown-signals.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS hwsec shutdown signals."
author          "chromium-os-dev@chromium.org"

start on stopping boot-services

task

script
    initctl emit hwsec-stop-clients-signal
    initctl emit hwsec-stop-tpm_manager-signal
    initctl emit hwsec-stop-low-level-tpm-daemon-signal
end script
