# swap_management/share/swap_start.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Setup compressed memory swap"
author          "ctshao@google.com"

start on started system-services
task

# Validate filesystem.
tmpfiles /usr/lib/tmpfiles.d/swap_management.conf

# It is exceedingly unlikely to OOM, as it's a simple dbus-send call,
# but better to kill it than to panic the system.
oom score -100

expect fork

exec minijail0 --config /usr/share/minijail/swap_management.conf \
    -- /usr/sbin/swap_management --swap_start
