# vm_tools/vhost_user_starter/init/vhost_user_starter.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Start the Vhost-user backend device"
author         "chromium-os-dev@chromium.org"

start on started vm_concierge
stop on stopped vm_concierge

# Upstart expects the main process to call fork(2) twice.
expect daemon

# If the job respawns 10 times in 10 seconds, stop trying.
respawn
respawn limit 10 10

# resourced handles low memory situations and shall not be OOM-killed.
oom score never

# Give time for vhost-user device to attempt an gracefully shutdown.
kill timeout 3

# The virtio-fs device opens a lot of fds and so the whole process tree needs a
# much higher limit.
limit nofile 1024 262144

script
# uid:
# 0 0 1 root
# 299 299 1 crosvm
# 600 600 50,65534 65534 1,655360 655360 2000000 kAndroidUidMap
uid_map="0 0 1,299 299 1,600 600 50,65534 65534 1,655360 655360 2000000"

# gid:
# 0 0 1 root
# 299 299 1 crosvm
# 600 600 50,65534 65534 1,655360 655360 2000000 kAndroidGidMap
gid_map="0 0 1,299 299 1,600 600 50,20119 20119 1,65534 65534 1,\
655360 655360 2000000"

# vhost_starter_daemon needs to enter concierge mount namespace for vhost-user
# virtio-fs accessing /run/arcvm/
exec nsenter --mount=/run/namespaces/mnt_concierge --no-fork \
    -- minijail0 \
        -Kslave \
        --config /usr/share/minijail/vhost_user_starter.conf \
    -- /sbin/minijail0 \
        -Kslave \
        -U -pv -I -i \
        -c 'cap_setuid,cap_setgid+eip' --ambient \
        -k 'proc,/proc,proc,MS_NOSUID|MS_NODEV|MS_NOEXEC' \
        -m"${uid_map}" \
        -M"${gid_map}" \
        -- /usr/bin/vhost_user_starter
end script
