# flexor/ramfs/etc/init/boot-services.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Job to trigger boot services"
author          "chromeos-flex-remote-deployment@google.com"

# Stub boot-services script.
start on stopped pre-startup

# Disable OOM killer as we must never fail.
oom score never
