# flexor/ramfs/etc/init/debug_tty.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Debug tty that only runs in debug builds"
author          "chromeos-flex-remote-deployment@google.com"

console output

# Starts a debug console. For now we always have one.
# TODO(b/314819758): Only attach a debug console if we are a debug build.
start on starting boot-services
respawn

# Disable OOM killer as we must never fail.
oom score never

exec /bin/sh
