# flexor/ramfs/etc/init/flexor.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Startup for Flexor script"
author          "chromeos-flex-remote-deployment@google.com"

# Flexor depends on udevd and syslog being fully initialized.
start on stopped udev-trigger-early and started syslog
expect fork

# Disable OOM killer as we must never fail.
oom score never

exec /usr/sbin/flexor
