# flexor/ramfs/etc/init/syslog.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Syslog daemon"
author        "chromeos-flex-remote-deployment@google.com"

# Starts the rsyslog daemon process inside minijail and runs after the
# boot-services have been started.
# This file is a fork of syslog.conf in CrOS.
start on starting boot-services
expect fork
respawn

# Disable OOM killer as we must never fail.
oom score never

exec minijail0 -l -N --uts -i -e -t -T static \
      /usr/sbin/rsyslogd -n -f /etc/rsyslog.chromeos -i /tmp/rsyslogd.pid
