# flexor/ramfs/etc/init/udev.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the udev daemon"
author          "chromeos-flex-remote-deployment@google.com"

# Starts the udev daemon process. The general "trigger"s for devices are in
# udev-trigger-* jobs.
# This file is a fork of udev.conf in CrOS.
start on stopped pre-startup
respawn

# Disable OOM killer as we must never fail.
oom score never

expect fork
exec udevd --daemon
