# bootsplash/init/upstart/bootsplash.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
description   "Signal Chromium OS bootsplash"
author        "chromium-os-dev@chromium.org"
# OK to kill if absolutely needed.
oom score -100

start on started dbus

script
  # Send SIGUSR1 to chargeplash, to indicate that DBus is up and running.
  if pgrep bootsplash; then pkill --signal SIGUSR1 bootsplash; fi
end script
