# init/upstart/libsegmentation.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Seed feature level"
author          "chromium-os-dev@chromium.org"

oom score never

# Start alongside chromeos_startup.
# Unless boot-splash has already query CBX level, at first reboot after update,
# will query gsctool and seed /run/libsegmenation cache for boards that can
# support CBX.
# For subsequent boot, the tool just read the VPD and confirms the data is
# valid.
# For boards that do not support CBX, this is mostly a noop operation.
start on stopped boot-splash

exec feature_check --feature_level
