# mojo_service_manager/init/mojo_service_manager.conf
# Copyright 2022 The ChromiumOS Authors.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chrome OS mojo service manager"
author          "chromium-os-dev@chromium.org"

# TODO(b/235922792): Restarting with `ui` is a workaround. We will remove this
# after migrating all the depended services to use service manager. This can
# only be done after we enable the service manager.
start on starting boot-services or starting ui
stop on stopped boot-services or stopping ui

respawn
# If the job respawns 3 times in 10 seconds, stop trying.
respawn limit 3 10

# -900 for CrOS daemons that are needed to auto-update, but can restart. Mojo is
# used by update_engine, so this is needed to auto-update.
oom score -900

tmpfiles /usr/lib/tmpfiles.d/on-demand/mojo_service_manager.conf

script
  # TODO(237029609): Move sandboxing logic into mojo_service_manager daemon.
  set --
  if [ -d /usr/local/etc/mojo/service_manager/policy ]; then
    set -- "$@" \
      -k \
      'tmpfs,/usr/local,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC,mode=755,size=10M' \
      -b /usr/local/etc/mojo/service_manager/policy
  fi

  exec minijail0 -u mojo-service-manager -g mojo-service-manager -G \
    --profile=minimalistic-mountns --uts -l -p -N \
    -k 'tmpfs,/run,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC,mode=755,size=10M' \
    -b /run/mojo,/run/mojo,1 \
    "$@" \
    -S /usr/share/policy/mojo_service_manager.policy \
    -- /usr/bin/mojo_service_manager
end script
