# pciguard/insecure-mode/pciguard-insecuremode.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on started system-services

description "Task to disable/enable pciguard on startup on test images"
author      "chromium-os-dev@chromium.org"

oom score -200

script

  if [ ! -f /var/lib/misc/pciguard_insecure_mode ]; then
    logger -p INFO "pciguard-insecuremode: normal mode detected, start pciguard"
    start pciguard
    exit 0
  fi

  exec /usr/share/cros/init/pciguard-insecuremode.sh

end script
