# minios/init/boot-services.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Job to trigger boot services"
author          "chromeos-core-services@google.com"

# A stub script to start `dbus.conf` which depends on boot-services.
start on started tmpfile

oom score never
