# minios/init/debug-tty.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Debug tty that only runs in debug builds"
author          "chromeos-core-services@google.com"

env DEBUG_CONSOLE=/dev/pts/2

console output

start on started frecon
respawn

oom score never

pre-start script
  if crossystem "cros_debug?1"; then
    return
  fi
  # stop job from continuing since we are not in dev mode
  stop; exit 0;
end script

exec /bin/sh<"${DEBUG_CONSOLE}" >"${DEBUG_CONSOLE}" 2>&1
