# minios/init/minios_util_logs.conf
# Copyright 2024 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Job to fetch NBR logs"
author          "chromeos-core-services@google.com"

start on started system-services

task

normal exit 0

# Killable for memory leaks.
oom score -100

tmpfiles /usr/lib/tmpfiles.d/on-demand/minios_util.conf

exec minijail0 --config /usr/share/minijail/minios_util.conf \
    -- /usr/bin/minios_util --retrieve=/var/log/minios \
    --erase --clear_key
