# minios/init/tmpfile.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start tmpfiles.d for D-Bus"
author          "chromeos-core-services@google.com"

# Brings up `systemd-tmpfiles` to create D-Bus directories before starting
# boot-services and dbus upstart jobs.
start on startup

oom score never

pre-start script
  systemd-tmpfiles --create --remove --boot \
     --prefix /dev --prefix /run \
     --prefix /home --prefix /var
end script
