# heartd/init/heartd.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "ChromeOS heartd (Health Ensure and Accident Resolve Treatment)"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

# It should not be killed by OOM killer, otherwise we may loose protection.
oom score never

# Run the daemon. It enters minijail in the daemon.
exec /usr/bin/heartd
