# ml_core/opencl_caching/init/opencl-cacher.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "OpenCL shader compilation cacher"
author "chromium-os-dev@chromium.org"

start on start-user-session
stop on stopping ui

oom score -100

exec /usr/bin/opencl_cacher
